# home/.config/upstart_local/dropbox.conf
description "User upstart for dropbox on the cli"
start on startup 
respawn
exec $HOME/.dropbox-dist/dropboxd
